# ansible/roles/run_backup/templates/borg-backup.timer
[Unit]
Description=Borg backup

[Timer]
OnUnitActiveSec=60

[Install]
WantedBy=timers.target

